# systemd/dc.target
[Unit]
Description=On DC power
DefaultDependencies=no
StopWhenUnneeded=yes
